# jobscheduler-joc/root/etc/services.d/joc/run
#!/usr/bin/with-contenv sh
#$SCHEDULER_HOME/scheduler/bin/jobscheduler.sh start
