# jobscheduler-joc/root/etc/services.d/joc/run
#!/usr/bin/with-contenv sh
$SCHEDULER_HOME/jetty/bin/jetty.sh run
