# pentaho-baserver/root/etc/services.d/baserver/run
#!/usr/bin/with-contenv sh

if [ -f $BASERVER_HOME/.baserver_configured ]; then
   cd $BASERVER_HOME
   sudo -u pentaho ./start-pentaho.sh
fi;
