# pentaho-pdi/root/etc/services.d/pdi-node/run
#!/bin/bash

# cd /usr/pentaho/data-integration
# sudo -H -u pentaho ./carte.sh ./config/carte-config.xml

